# systemd/llama-server.service
[Unit]
Description=llama.cpp OpenAI-compatible API server
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
# Replace these via 10-install-systemd.sh — placeholders are filled in there.
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=MODEL=__MODEL__
Environment=LLAMA_DIR=__LLAMA_DIR__
Environment=LLAMA_BIN=__LLAMA_DIR__/llama-server
Environment=LLAMA_PORT=1234
# LLAMA_BIND left blank so start-llama-server.sh picks the current LAN IP at
# boot — makes DHCP-renewed addresses Just Work.
Environment=LLAMA_NGL=99
# LLAMA_CTX left blank → auto-sized to the largest that fits in free VRAM.
ExecStart=__INSTALL_DIR__/start-llama-server.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
