# systemd/llama-server.service
[Unit]
Description=llama.cpp OpenAI-compatible API server
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
# Replace these via 10-install-systemd.sh — placeholders are filled in there.
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=MODEL=__MODEL__
Environment=LLAMA_DIR=__LLAMA_DIR__
Environment=LLAMA_BIN=__LLAMA_DIR__/llama-server
Environment=LLAMA_PORT=1234
# LLAMA_BIND left blank → start-llama-server.sh defaults to 0.0.0.0 so
# container port-mapping + DHCP-renewed addresses Just Work. Pin a specific
# interface with Environment=LLAMA_BIND=<ip>.
Environment=LLAMA_NGL=99
# LLAMA_CTX left blank → auto-sized to the largest that fits in free VRAM.
ExecStart=__INSTALL_DIR__/start-llama-server.sh
# A `systemctl stop` sends SIGTERM, and llama-server exits 143 (128+15).
# Without this, Restart=on-failure treats 143 as a failure and brings the
# unit right back after an explicit stop — breaking the dashboard's Unload
# button (VRAM appears to free, then fills again 5s later).
SuccessExitStatus=0 143 SIGTERM
Restart=on-failure
RestartSec=5
# Make sure SIGTERM is what actually goes to the process, and wait long
# enough for the model KV cache to be freed before escalating to SIGKILL.
KillSignal=SIGTERM
TimeoutStopSec=30
# Remove memlock cap so --mlock can pin the model weights without hitting
# the default 8 MiB RLIMIT_MEMLOCK. Harmless if --mlock is off.
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
