# systemd/lmstudio-dashboard.service
[Unit]
Description=ML Stack Trainer Dashboard
After=network-online.target
Wants=network-online.target
# Safety net: at most 5 failed starts in 5 minutes, then systemd gives up.
# Prevents runaway restart loops (we once hit 107 restarts) if the venv is
# missing or a dep is busted. These keys live in [Unit] for compat with
# older systemd — newer systemd also accepts them in [Service].
# Run `sudo systemctl reset-failed lmstudio-dashboard.service` to retry.
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
# Replace these via 10-install-systemd.sh — placeholders are filled in there.
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=INSTALL_DIR=__INSTALL_DIR__
Environment=AUTO_LAUNCH_LMSTUDIO=0
# DASHBOARD_HOST intentionally unset → start-dashboard picks the current
# LAN IP at boot. Override with an explicit IP to pin the bind.
Environment=DASHBOARD_PORT=8765
Environment=VENV_DIR=__VENV_DIR__
Environment=LMSTUDIO_DIR=__LMSTUDIO_DIR__
Environment=LLAMA_DIR=__LLAMA_DIR__
Environment=LLAMA_BIN=__LLAMA_DIR__/llama-server
Environment=GGUF_MODELS_DIR=__HOME__/models
ExecStart=__INSTALL_DIR__/09-start-dashboard.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
